# benchmarks/retry-storm-duplicate-transfer/app/supervisord.conf
[supervisord]
nodaemon=false
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[program:processor]
command=python3 -u /workspace/app/processor.py
autostart=true
autorestart=true
stdout_logfile=/tmp/processor.log
stderr_logfile=/tmp/processor_err.log
priority=10

[program:chaos_proxy]
command=python3 -u /workspace/app/chaos_proxy.py
autostart=true
autorestart=true
stdout_logfile=/tmp/chaos_proxy.log
stderr_logfile=/tmp/chaos_proxy_err.log
priority=20

[program:gateway]
command=python3 -u /workspace/app/gateway.py
autostart=true
autorestart=true
stdout_logfile=/tmp/gateway.log
stderr_logfile=/tmp/gateway_err.log
priority=30
